# app-supervisor.conf
[supervisord]
nodaemon = true

[program: springboot-scala-example]
command = /bin/bash -c "/usr/local/app/bin/run-app.sh"
stdout_logfile = /usr/local/app/logs/stdout.log
user = app
